# bitfocus-companion/rootfs/etc/services.d/companion/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take action on shutdown
# ==============================================================================

bashio::log.info "Bitfocus Companion service is gestopt."
