# bitfocus-companion/rootfs/etc/services.d/companion/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Uitvoeren als de Companion service is gestopt
# ==============================================================================

bashio::log.info "Bitfocus Companion is gestopt"
