# bitfocus-companion/rootfs/etc/services.d/companion/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Bitfocus Companion service
# ==============================================================================

bashio::log.info "Starten van Bitfocus Companion..."

# Zorg dat de configuratiedirectory is ingesteld
export COMPANION_CONFIG_BASEDIR="/data/companion"

# Start Companion met de juiste configuratie
cd /app
exec /docker-entrypoint.sh sh -c "./node-runtimes/main/bin/node ./main.js --admin-address 0.0.0.0 --admin-port 8000 --config-dir ${COMPANION_CONFIG_BASEDIR} --extra-module-path /app/module-local-dev"
