# bitfocus-companion/rootfs/etc/services.d/companion/run
#!/bin/bash
set -e

echo "Bitfocus Companion proxy addon starting..."

# Maak data directory voor persistentie
echo "Creating persistent storage directory..."
mkdir -p /data/companion
chmod 777 /data/companion

# Bepaal architectuur
ARCH=$(uname -m)
echo "Detected architecture: ${ARCH}"

# Start socat proxies voor de poorten
echo "Starting port forwarding..."
socat TCP-LISTEN:8000,fork TCP:127.0.0.1:8000 &
socat TCP-LISTEN:16622,fork TCP:127.0.0.1:16622 &
socat TCP-LISTEN:16623,fork TCP:127.0.0.1:16623 &
socat TCP-LISTEN:28492,fork TCP:127.0.0.1:28492 &

# Download en installeer Companion direct
echo "Downloading and installing Companion..."
cd /tmp
if [ "${ARCH}" = "aarch64" ]; then
  COMPANION_ARCH="arm64"
else
  COMPANION_ARCH="x64"
fi

echo "Using Companion architecture: ${COMPANION_ARCH}"
curl -L -o companion.tar.gz https://bitfocus.io/companion/download/linux-${COMPANION_ARCH}
mkdir -p /opt/companion
tar -xzf companion.tar.gz -C /opt/companion --strip-components=1
cd /opt/companion

# Start Companion
echo "Starting Companion..."
node main.js --admin-address 0.0.0.0 --admin-port 8000 --config-dir /data/companion &

# Houd de addon actief
echo "Addon running, monitoring Companion process..."
while true; do
  sleep 30
  # Controleer of het node proces nog draait
  if ! pgrep -f "node main.js" > /dev/null; then
    echo "Companion process stopped, restarting..."
    cd /opt/companion
    node main.js --admin-address 0.0.0.0 --admin-port 8000 --config-dir /data/companion &
  fi
done
